# files/package/base-files/files/etc/hotplug.d/iface/99-netfix
#!/bin/sh
case "$ACTION" in
ifup)
/sbin/video.sh
;;
ifdown)
;;
esac
